# embed-fixer.service
[Unit]
Description=Discord Embed Fixer Bot
After=network-online.target
Wants=network-online.target
Restart=always

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/embed-fixer
Environment=PATH=/home/admin/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/admin/embed-fixer/.env
ExecStart=/home/admin/.local/bin/uv run main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=embed-fixer

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/admin/embed-fixer
ReadWritePaths=/home/admin/.cache

[Install]
WantedBy=multi-user.target

